# etc/runit/sv/gopreload/run
#!/bin/sh
exec /usr/share/gopreload/bin/StartPreload.sh 2>&1